# deploy/pdf-checker-worker.service
[Unit]
Description=PDF external-source audit worker
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/pdf-checker/current
Environment=NODE_ENV=production
EnvironmentFile=/etc/pdf-checker.env
ExecStart=/usr/local/bin/node /opt/pdf-checker/current/dist/audit-worker.mjs
Restart=on-failure
RestartSec=5
TimeoutStopSec=30
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/pdf-checker

[Install]
WantedBy=multi-user.target
